# nginx/supervisor_impianti_termici.conf
[program:gtfs-server]
#this would help with env vars...
#command = source /var/opt/impianti_termici/environment.sh &&
command = /usr/local/bin/uwsgi -i /var/opt/gtfs-server/nginx/uwsgi.ini
